# web/config/supervisord.conf
[program:php7.3-fpm]
command=service php7.4-fpm start

[program:nginx]
command=service nginx start

[program:redis]
command=service redis-server start
